# uwsgi.service
[Unit]
Description=uWSGI Service

[Service]
ExecStart=chagi-env/bin/uwsgi --ini hps-django/django-beta/mybend_uwsgi.ini
